# scripts/launchd/com.cortex.digest-cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cortex.digest-cron</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>source ~/.cortex/env &amp;&amp; curl -sf -X POST -H "Authorization: Bearer $CORTEX_API_TOKEN" http://127.0.0.1:8420/api/v1/digest/push?days=1 || true</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/cortex/digest-cron.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/cortex/digest-cron.err</string>
</dict>
</plist>
