# scripts/launchd/com.cortex.ilink-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cortex.ilink-agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/cuiliangjing/Projects/cortex-wechat/scripts/launchd/ilink-agent.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/cuiliangjing/Library/Logs/cortex-wechat/ilink-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/cuiliangjing/Library/Logs/cortex-wechat/ilink-agent.err</string>
</dict>
</plist>
